# hmi-x.service
[Unit]
Description=HMI X session via xinit (for CLI-only systems)
After=network.target

[Service]
Type=simple
User=pi
Group=pi
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStart=/usr/bin/xinit /opt/hmi/venv/bin/python3 /opt/hmi/main.py -- :0 vt1 -nolisten tcp
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
